# openstack/acmeair/acmeair_dockerfile/acmeair_webserver/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:liberty]
command=/acmeair/wlp/bin/server start server1
